# omni/libvirt-provider/omni-infra-provider-libvirt.service
[Unit]
Description=Omni infrastructure provider for libvirt
Documentation=https://github.com/siderolabs/omni-infra-provider-libvirt
Wants=network-online.target
After=network-online.target libvirtd.service
Requires=libvirtd.service

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=/etc/omni-infra-provider-libvirt/provider.env
ExecStart=/usr/local/bin/omni-infra-provider-libvirt --config-file=/etc/omni-infra-provider-libvirt/config.yaml --id=libvirt --provider-name="Libvirt vanillax-dell" --provider-description="KVM host on CachyOS 192.168.10.20" --image-cache-path=/var/cache/omni-infra-provider-libvirt
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/var/cache/omni-infra-provider-libvirt

[Install]
WantedBy=multi-user.target
