# omni/libvirt-provider/talos-wired-pod-routes.service
# INTERIM: static return routes to the wired nodes' PodCIDRs. Replaced by
# FRR + Cilium BGP in the follow-up; delete this unit when FRR owns these.
[Unit]
Description=Interim return routes for wired Talos PodCIDRs
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/ip route replace 10.244.0.0/24 via 192.168.10.177 dev wlan0
ExecStart=/usr/bin/ip route replace 10.244.1.0/24 via 192.168.10.81 dev wlan0
ExecStop=-/usr/bin/ip route del 10.244.0.0/24 via 192.168.10.177 dev wlan0
ExecStop=-/usr/bin/ip route del 10.244.1.0/24 via 192.168.10.81 dev wlan0

[Install]
WantedBy=multi-user.target
