# omni/libvirt-provider/talos-routed-egress-nat.service
[Unit]
Description=Public egress NAT for routed Talos libvirt guests
Wants=network-online.target
After=network-online.target ufw.service libvirtd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/libexec/talos-routed-egress-nat start
ExecStop=/usr/local/libexec/talos-routed-egress-nat stop

[Install]
WantedBy=multi-user.target
